# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[program:queue-worker]
command=php /var/www/html/artisan queue:work --queue=emails,default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
; Supervisor tidak menyetel HOME saat drop privilege lewat `user=`, jadi tanpa
; ini worker mewarisi HOME=/root yang tidak bisa ditulis www-data.
environment=HOME="/var/www",USER="www-data"
; Kegagalan saat boot biasanya transient (Postgres eksternal belum reachable).
; Dengan startretries kecil Supervisor cepat menyerah dan menaruh worker di state
; FATAL yang TIDAK pernah di-retry lagi -> harus di-start manual.
startretries=100
startsecs=5
; Docker hanya memberi stop_grace_period (default 10s) sebelum SIGKILL, jadi nilai
; yang lebih besar dari itu percuma. Disamakan dengan docker-compose.yml.
stopwaitsecs=120
stopsignal=TERM
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:scheduler]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
environment=HOME="/var/www",USER="www-data"
startretries=100
startsecs=5
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
